# blackbox/lxqt-kwin-dark/desktop-overlay/etc/skel/.config/autostart/kdeconnect-indicator.desktop
[Desktop Entry]
Exec=kdeconnect-indicator
Name=kdeconnect-indicator
Type=Application
Icon=kdeconnect-indicator
Version=1.0
X-LXQt-Need-Tray=true
