# oem/starlabs-oem-kde/live-overlay/etc/skel/.config/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Install Manjaro Linux
GenericName=Live Installer
Comment=Install the operating system to disk
Comment[de]=Manjaro Linux installieren
X-GNOME-Autostart-enabled=true
Exec=kstart5 --fullscreen calamares_polkit
Icon=calamares
Terminal=false
StartupNotify=true
Type=Application
Categories=Qt;System;
